# deploy/ai.mandel.casino-review.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.mandel.casino-review</string>

    <!-- The wrapper loads .env then execs the prebuilt binary. -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/kalata/MandelAI/testing/casino-review/scripts/casino-review.sh</string>
    </array>

    <!-- Start at login and keep it alive: relaunch if it ever exits/crashes. -->
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>

    <!-- Don't hammer relaunches if it's crash-looping. -->
    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/Users/kalata/Library/Logs/casino-review.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kalata/Library/Logs/casino-review.log</string>
</dict>
</plist>
